# deploy/components/cron/crontab
* * * * * php /code/artisan schedule:run